# screenshot-bot/systemd/screenshot-bot.service
[Unit]
Description=Screenshot Telegram Bot
After=network.target

[Service]
Type=simple
User=tabris
WorkingDirectory=/home/tabris/bots/screenshot-bot
ExecStart=/home/tabris/bots/screenshot-bot/venv/bin/python3 bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
