# screenshot-bot/systemd/screenshot-cleanup.timer
[Unit]
Description=Run screenshot cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
